# deploy/capsule-video-worker.service
[Unit]
Description=Capsule video import worker (yt-dlp / ffmpeg)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/capsule
Environment=APP_ENV=prod
Environment=DB_DSN=sqlite:/var/www/capsule/data/database.sqlite
Environment=VIDEO_IMPORT_ROOT=/var/www/capsule/public/uploads/media/imports
ExecStart=/usr/bin/php /var/www/capsule/bin/video-import-worker.php
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
